# root/etc/services.d/x11vnc/run
#!/usr/bin/with-contenv sh
 env
/usr/bin/x11vnc -display ${DISPLAY} -nopw -wait 5 -shared
