# root/etc/services.d/xvfb/run
#!/usr/bin/with-contenv sh

/usr/bin/Xvfb ${DISPLAY} -screen 0 ${ENV_RESOLUTION:-1920x1080x24}
